# profile.d/dev.sh
#!/usr/bin/env bash

#source "${ENV_LIB}/dev.sh"

app "${HOME}/dev"

